# profile.d/mark-symlinked-directories.sh
#!/bin/sh
bind 'set mark-symlinked-directories on'
