# profile.d/README
Put these in /etc/profile.d, perm 644, root:root
